# .config/autostart/udiskie.desktop
[Desktop Entry]
Name=Udiskie
GenericName=Udiskie
Exec=udiskie
Terminal=false
Type=Application
Categories=
StartupNotify=false
