# .config/autostart/xsettingsd.desktop
[Desktop Entry]
Name=XSettingsd
GenericName=Settings daemon
Exec=xsettingsd
Terminal=false
Type=Application
Categories=
StartupNotify=false
